# ub16/supervisord.conf
[supervisord]
nodaemon=true

[program:tomcat7]
command=/tomcat/bin/catalina.sh run
stdout_logfile=/var/log/supervisor/%(program_name)s.log
redirect_stderr = true
autorestart=true